# scripts/nova-openauto.service
[Unit]
Description=OpenAuto Android Auto
After=nova-weston.service avahi-daemon.service
Requires=nova-weston.service
Wants=avahi-daemon.service

[Service]
Type=simple
User=root
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=WAYLAND_DISPLAY=wayland-1
Environment=QT_QPA_PLATFORM=wayland
Environment=QT_WAYLAND_DISABLE_WINDOWDECORATION=1
ExecStart=/opt/openauto/bin/autoapp --config /etc/openauto/openauto.ini
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
